# src/main/resources/cmd/s2si.service
[Unit]
Description=S2SI
After=syslog.target

[Service]
User=<user>
ExecStart=<path>
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target